# .trash/_remote/launchd/com.chinapolicymonitor.macro.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chinapolicymonitor.macro</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/sd/miniforge3/bin/python3</string>
        <string>/Users/sd/PycharmProjects/CMM/run_fetch_macro.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sd/PycharmProjects/CMM</string>

    <!--
        Triggered once per day at 03:00 local time.
        The script itself sleeps a random 0-18h offset before fetching,
        so actual fetch time varies across the day.
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sd/PycharmProjects/CMM/data/logs/launchd_macro_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sd/PycharmProjects/CMM/data/logs/launchd_macro_stderr.log</string>
</dict>
</plist>
