# .trash/_remote/launchd/com.chinapolicymonitor.news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chinapolicymonitor.news</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/sd/miniforge3/bin/python3</string>
        <string>/Users/sd/PycharmProjects/CMM/run_fetch_news.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sd/PycharmProjects/CMM</string>

    <!-- Every hour -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/sd/PycharmProjects/CMM/data/logs/launchd_news_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sd/PycharmProjects/CMM/data/logs/launchd_news_stderr.log</string>
</dict>
</plist>
