# backend/scheduler/com.chinapolicymonitor.policies.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chinapolicymonitor.policies</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cmm-fetch</string>
        <string>policies</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/main/PycharmProjects/CMM</string>

    <!-- Every 12 hours -->
    <key>StartInterval</key>
    <integer>43200</integer>

    <key>StandardOutPath</key>
    <string>/Users/main/PycharmProjects/CMM/data/logs/launchd_policies_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/main/PycharmProjects/CMM/data/logs/launchd_policies_stderr.log</string>
</dict>
</plist>
